# init_scripts/configurator.conf
description "Configurator server"
author "upendra.sahu@maplelabs.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /home/user/WORK/configurator-exporter/api-server.py
